# tools/launchd/com.scottdube.ha.energy-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LRD whole-home energy daily audit.

  Runs energy_audit.py once per day at 06:30 EDT against yesterday's
  local-day data. Yesterday is fully recorded by then (midnight rollover +
  margin for late publishes from ESPHome Vue panels), and 06:30 is early
  enough that anomaly notifications reach Scott before the morning routine.

  Designed for the always-on Mac mini at LRD. Reuses the existing
  ~/.venv/zwave-health venv (websockets already installed).

  Install:
      ln -sf /Users/scottdube/code/home-assistant/tools/launchd/com.scottdube.ha.energy-audit.plist \
             ~/Library/LaunchAgents/
      launchctl bootstrap gui/$(id -u) \
             ~/Library/LaunchAgents/com.scottdube.ha.energy-audit.plist

  Force a run immediately (e.g. after token rotation):
      launchctl kickstart -k gui/$(id -u)/com.scottdube.ha.energy-audit

  Tail the log:
      tail -f ~/Library/Logs/ha-energy-audit.log

  Uninstall:
      launchctl bootout gui/$(id -u) \
             ~/Library/LaunchAgents/com.scottdube.ha.energy-audit.plist
      rm ~/Library/LaunchAgents/com.scottdube.ha.energy-audit.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scottdube.ha.energy-audit</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/scottdube/code/home-assistant &amp;&amp; git pull --ff-only --quiet 2&gt;/dev/null; /Users/scottdube/.venv/zwave-health/bin/python3 /Users/scottdube/code/home-assistant/tools/energy_audit.py --print-clean</string>
    </array>

    <!-- 06:30 LRD-local daily. macOS launchd interprets StartCalendarInterval in
         the machine's local time, so this fires at 06:30 EDT when the Mac mini
         is at LRD. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/scottdube/Library/Logs/ha-energy-audit.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/scottdube/Library/Logs/ha-energy-audit.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
