# tools/launchd/com.scottdube.ha.zwave-health-probe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Z-Wave node health probe — node 55 (Fibaro FGD-212 lamp post dimmer at LRD).

  Runs zwave_health_probe.py every 4 hours, appending one row per run to
  ~/zwave-health/node-55.csv. RunAtLoad fires the probe immediately on
  bootstrap.

  This is the canonical reference instance of the Z-Wave node health
  monitoring pattern described in ADR-021. Other monitored nodes have
  their own sibling plists under this directory.

  Designed for the always-on Mac mini at LRD.

  Install:
      python3 -m venv ~/.venv/zwave-health
      ~/.venv/zwave-health/bin/pip install websockets
      ln -sf /Users/scottdube/code/home-assistant/tools/launchd/com.scottdube.ha.zwave-health-probe.plist \
             ~/Library/LaunchAgents/
      launchctl bootstrap gui/$(id -u) \
          ~/Library/LaunchAgents/com.scottdube.ha.zwave-health-probe.plist

  Tail the log:
      tail -f ~/Library/Logs/zwave-health-probe-node-55.log

  Uninstall:
      launchctl bootout gui/$(id -u) \
          ~/Library/LaunchAgents/com.scottdube.ha.zwave-health-probe.plist
      rm ~/Library/LaunchAgents/com.scottdube.ha.zwave-health-probe.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scottdube.ha.zwave-health-probe-node-55</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/scottdube/.venv/zwave-health/bin/python3</string>
        <string>/Users/scottdube/code/home-assistant/tools/zwave_health_probe.py</string>
        <string>--node</string>
        <string>55</string>
        <string>--rounds</string>
        <string>5</string>
        <string>--csv</string>
        <string>/Users/scottdube/zwave-health/node-55.csv</string>
    </array>

    <key>StartInterval</key>
    <integer>14400</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/scottdube/Library/Logs/zwave-health-probe-node-55.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/scottdube/Library/Logs/zwave-health-probe-node-55.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
