# pool/scripts/launchd/com.scottdube.ha.pool-log-rsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Pulls /config/pool_state_log.csv from the HA NUC every 5 minutes via rsync
  over SSH, writing into the local repo's pool/analysis/ directory.

  Install:
      cp pool/scripts/launchd/com.scottdube.ha.pool-log-rsync.plist \
         ~/Library/LaunchAgents/
      launchctl load ~/Library/LaunchAgents/com.scottdube.ha.pool-log-rsync.plist

  Adjust port (22 vs 2222) before loading if your HA SSH addon uses 2222.

  See pool/docs/rsync-setup.md for full setup instructions and troubleshooting.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scottdube.ha.pool-log-rsync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/rsync</string>
        <string>-az</string>
        <string>--timeout=30</string>
        <string>-e</string>
        <string>ssh -p 22 -i /Users/scottdube/.ssh/ha_nuc_rsync -o StrictHostKeyChecking=accept-new -o ConnectTimeout=10</string>
        <string>sdube@192.168.50.11:/config/pool_state_log.csv</string>
        <string>/Users/scottdube/code/home-assistant/pool/analysis/pool_state_log_live.csv</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/scottdube/Library/Logs/ha-pool-rsync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/scottdube/Library/Logs/ha-pool-rsync.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
